# scripts/com.thebrinsider.receita-qsa-mensal.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.thebrinsider.receita-qsa-mensal</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/luizlessa/transparencia-federal/scripts/receita-qsa-cron.sh</string>
  </array>

  <!-- Dia 15 de cada mês às 09:00 BRT (Receita publica no início do mês) -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Day</key>
    <integer>15</integer>
    <key>Hour</key>
    <integer>9</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/Users/luizlessa/Library/Logs/receita-qsa-cron.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/luizlessa/Library/Logs/receita-qsa-cron-error.log</string>

  <!-- QSA demora ~10 min; throttle de 1 dia evita reexecuções acidentais -->
  <key>ThrottleInterval</key>
  <integer>86400</integer>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/opt/node@22/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>HOME</key>
    <string>/Users/luizlessa</string>
  </dict>
</dict>
</plist>
