# scripts/com.thebrinsider.camara-mensal.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.thebrinsider.camara-mensal</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/luizlessa/transparencia-federal/scripts/camara-cron.sh</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/opt/node@22/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>HOME</key>
    <string>/Users/luizlessa</string>
  </dict>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Day</key>
    <integer>6</integer>
    <key>Hour</key>
    <integer>10</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/luizlessa/Library/Logs/camara-cron.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/luizlessa/Library/Logs/camara-cron.log</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
